# backend/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:uwsgi]
directory=/app/conf
command=/usr/local/bin/uwsgi --ini uwsgi.ini
