# frontend/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=rc-service nginx start
